# etc/systemd/system/customwifi.service
[Unit]
Description=Custom wifi unit

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/customwifi

[Install]
